# app-admin/chef/files/initd/chef
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
	ebegin "Starting Chef client"
	start-stop-daemon --start --quiet -b -p /var/run/chef/chef-client.pid --exec /usr/bin/chef-client -- -c /etc/chef/client.rb
	eend $?
}

stop() {
	ebegin "Shutting down Chef client"
	start-stop-daemon --stop --quiet -p /var/run/chef/chef-client.pid
	eend $?
}
